# supervisor.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /usr/src/app/uwsgi.ini --logto /tmp/uwsgi.log

[program:nginx-app]
command = /usr/sbin/nginx
